# odoo/rootfs/etc/s6-overlay/s6-rc.d/init-nginx/run
#!/command/with-contenv bash
# shellcheck shell=bash
set -e

echo "Configuring Nginx reverse proxy for Odoo..."

# Generate nginx config from template
# WOOWTECH: Always use HTTP for LAN access
printf '{
  "listen_port": "8069",
  "backend_port": "8169",
  "websocket_port": "8072",
  "ssl": false
}' | tempio \
    -template /etc/nginx/templates/odoo.gtpl \
    -out /etc/nginx/servers/odoo.conf

echo "Nginx configuration generated."
