# rootfs/etc/s6-overlay/s6-rc.d/svc-odoo/run
#!/command/with-contenv bash
# shellcheck shell=bash
set -e

# Wait for PostgreSQL to be ready
echo "Waiting for PostgreSQL..."
for i in $(seq 1 60); do
    if pg_isready -h localhost -p 5432 -U odoo >/dev/null 2>&1; then
        echo "PostgreSQL is ready."
        break
    fi
    if [ "$i" -eq 60 ]; then
        echo "ERROR: PostgreSQL did not become ready in time!"
        exit 1
    fi
    sleep 2
done

OPTIONS_SOURCE=${ADDON_DATA_FOLDER:-"/data"}/options.json
DB_NAME=$(jq -r '.odoo_db_name // "odoo"' "$OPTIONS_SOURCE")

# Check if database needs initialization (base module not installed)
DB_HAS_MODULES=$(s6-setuidgid postgres psql -U odoo -d "$DB_NAME" -tAc "SELECT COUNT(*) FROM pg_tables WHERE schemaname='public' AND tablename='ir_module_module';" 2>/dev/null || echo "0")

if [ "$DB_HAS_MODULES" = "0" ] || [ -z "$DB_HAS_MODULES" ]; then
    echo "Initializing Odoo database (first run, installing base module)..."
    s6-setuidgid odoo odoo \
        --config=/etc/odoo/odoo.conf \
        -i base \
        -d "$DB_NAME" \
        --stop-after-init \
        --without-demo=all
    echo "Odoo database initialization complete."
fi

echo "Starting Odoo 18 server..."

# When workers > 0, Odoo automatically spawns HTTP workers + gevent worker
# for WebSocket/longpolling on gevent_port (8072)
exec s6-setuidgid odoo odoo \
    --config=/etc/odoo/odoo.conf
