# odoo/rootfs/etc/s6-overlay/s6-rc.d/init-odoo-config/run
#!/command/with-contenv bash
# shellcheck shell=bash
set -e

OPTIONS_SOURCE=${ADDON_DATA_FOLDER:-"/data"}/options.json

echo "Generating Odoo configuration..."

# Read options
DB_PASSWORD=$(jq -r '.db_password // "odoo"' "$OPTIONS_SOURCE")
ADMIN_PASSWORD=$(jq -r '.admin_password // empty' "$OPTIONS_SOURCE")
WORKERS=$(jq -r '.odoo_workers // 2' "$OPTIONS_SOURCE")
MAX_CRON=$(jq -r '.odoo_max_cron_threads // 1' "$OPTIONS_SOURCE")
MEM_HARD=$(jq -r '.odoo_limit_memory_hard // 2684354560' "$OPTIONS_SOURCE")
MEM_SOFT=$(jq -r '.odoo_limit_memory_soft // 2147483648' "$OPTIONS_SOURCE")
TIME_CPU=$(jq -r '.odoo_limit_time_cpu // 60' "$OPTIONS_SOURCE")
TIME_REAL=$(jq -r '.odoo_limit_time_real // 120' "$OPTIONS_SOURCE")
LOG_LEVEL_RAW=$(jq -r '.odoo_log_level // "info"' "$OPTIONS_SOURCE")
# Map to Python logging uppercase levels for log_handler
case "$LOG_LEVEL_RAW" in
    debug) LOG_HANDLER_LEVEL="DEBUG" ;;
    info)  LOG_HANDLER_LEVEL="INFO" ;;
    warn)  LOG_HANDLER_LEVEL="WARNING" ;;
    error) LOG_HANDLER_LEVEL="ERROR" ;;
    critical) LOG_HANDLER_LEVEL="CRITICAL" ;;
    *)     LOG_HANDLER_LEVEL="INFO" ;;
esac
PROXY_MODE=$(jq -r '.odoo_proxy_mode // true' "$OPTIONS_SOURCE")
EXTRA_ADDONS=$(jq -r '.odoo_extra_addons // empty' "$OPTIONS_SOURCE")
DB_NAME=$(jq -r '.odoo_db_name // "odoo"' "$OPTIONS_SOURCE")
WITHOUT_DEMO=$(jq -r '.odoo_without_demo // true' "$OPTIONS_SOURCE")

# Build addons path
ADDONS_PATH="/mnt/extra-addons,/usr/lib/python3/dist-packages/odoo/addons"
if [ -n "$EXTRA_ADDONS" ]; then
    ADDONS_PATH="${EXTRA_ADDONS},${ADDONS_PATH}"
fi

# Ensure Odoo data directory exists
mkdir -p /data/odoo /data/logs/odoo
chown -R odoo:odoo /data/odoo /data/logs/odoo

# Generate odoo.conf
cat > /etc/odoo/odoo.conf <<EOF
[options]
; Database
db_host = localhost
db_port = 5432
db_user = odoo
db_password = ${DB_PASSWORD}
db_name = ${DB_NAME}
db_maxconn = 64

; Paths
data_dir = /data/odoo
addons_path = ${ADDONS_PATH}

; Server
http_port = 8169
gevent_port = 8072
proxy_mode = ${PROXY_MODE}
server_wide_modules = web

; Workers
workers = ${WORKERS}
max_cron_threads = ${MAX_CRON}

; Limits
limit_memory_hard = ${MEM_HARD}
limit_memory_soft = ${MEM_SOFT}
limit_time_cpu = ${TIME_CPU}
limit_time_real = ${TIME_REAL}

; Logging
logfile = /data/logs/odoo/odoo.log
log_handler = :${LOG_HANDLER_LEVEL}

; Security
list_db = True
EOF

# Set admin password if provided
if [ -n "$ADMIN_PASSWORD" ]; then
    echo "admin_passwd = ${ADMIN_PASSWORD}" >> /etc/odoo/odoo.conf
fi

# Set without_demo
if [ "$WITHOUT_DEMO" = "true" ]; then
    echo "without_demo = all" >> /etc/odoo/odoo.conf
fi

chown odoo:odoo /etc/odoo/odoo.conf
chmod 640 /etc/odoo/odoo.conf

echo "Odoo configuration generated at /etc/odoo/odoo.conf"
