# odoo/rootfs/etc/s6-overlay/s6-rc.d/init-postgres/run
#!/command/with-contenv bash
# shellcheck shell=bash
set -e

OPTIONS_SOURCE=${ADDON_DATA_FOLDER:-"/data"}/options.json

# Read password from options
DB_PASSWORD=$(jq -r '.db_password // "odoo"' "$OPTIONS_SOURCE")
TZ_VAL=$(jq -r '.TZ // "Etc/UTC"' "$OPTIONS_SOURCE")

export POSTGRES_PASSWORD="$DB_PASSWORD"
export POSTGRES_USER="odoo"
export POSTGRES_DB="odoo"

echo "Initializing PostgreSQL ${PG_MAJOR}..."

# Ensure directories exist with correct ownership
mkdir -p "$PGDATA" /var/run/postgresql /data/logs/postgres
chown -R postgres:postgres "$PGDATA" /var/run/postgresql /data/logs/postgres
chmod 700 "$PGDATA" || true
chmod 3777 /var/run/postgresql || true

# Check if database is already initialized
if [ -s "$PGDATA/PG_VERSION" ]; then
    echo "PostgreSQL database directory already exists, skipping initialization."

    # Update password if changed
    pg_ctl -D "$PGDATA" -o "-c listen_addresses='localhost'" -w start -l /data/logs/postgres/init.log
    psql -U postgres -c "ALTER USER ${POSTGRES_USER} PASSWORD '${POSTGRES_PASSWORD}';" 2>/dev/null || true
    pg_ctl -D "$PGDATA" -m fast -w stop
else
    echo "Initializing new PostgreSQL database..."

    # Use nss_wrapper if needed
    local_uid="$(id -u)"
    if ! getent passwd "$local_uid" &>/dev/null; then
        for wrapper in {/usr,}/lib{/*,}/libnss_wrapper.so; do
            if [ -s "$wrapper" ]; then
                NSS_WRAPPER_PASSWD="$(mktemp)"
                NSS_WRAPPER_GROUP="$(mktemp)"
                export LD_PRELOAD="$wrapper" NSS_WRAPPER_PASSWD NSS_WRAPPER_GROUP
                printf 'postgres:x:%s:%s:PostgreSQL:%s:/bin/false\n' "$(id -u postgres)" "$(id -g postgres)" "$PGDATA" > "$NSS_WRAPPER_PASSWD"
                printf 'postgres:x:%s:\n' "$(id -g postgres)" > "$NSS_WRAPPER_GROUP"
                break
            fi
        done
    fi

    # Initialize the database as postgres user
    # Use a temp file for password instead of process substitution
    # (process substitution /dev/fd/N is not accessible after s6-setuidgid)
    PW_FILE="$(mktemp)"
    printf "%s\n" "$POSTGRES_PASSWORD" > "$PW_FILE"
    chmod 644 "$PW_FILE"

    s6-setuidgid postgres initdb \
        --username="$POSTGRES_USER" \
        --pwfile="$PW_FILE" \
        --encoding=UTF8 \
        --locale=en_US.UTF-8 \
        -D "$PGDATA"

    rm -f "$PW_FILE"

    # Cleanup nss_wrapper
    if [[ "${LD_PRELOAD:-}" == */libnss_wrapper.so ]]; then
        rm -f "$NSS_WRAPPER_PASSWD" "$NSS_WRAPPER_GROUP"
        unset LD_PRELOAD NSS_WRAPPER_PASSWD NSS_WRAPPER_GROUP
    fi

    # Configure pg_hba.conf for local connections
    {
        printf '\n# Allow local connections\n'
        printf 'local all all trust\n'
        printf 'host all all 127.0.0.1/32 md5\n'
        printf 'host all all ::1/128 md5\n'
    } >> "$PGDATA/pg_hba.conf"

    # Configure postgresql.conf
    {
        printf "\n# WOOWTECH Custom Settings\n"
        printf "listen_addresses = 'localhost'\n"
        printf "port = 5432\n"
        printf "max_connections = 100\n"
        printf "shared_buffers = 128MB\n"
        printf "effective_cache_size = 256MB\n"
        printf "work_mem = 4MB\n"
        printf "maintenance_work_mem = 64MB\n"
        printf "timezone = '%s'\n" "$TZ_VAL"
        printf "log_timezone = '%s'\n" "$TZ_VAL"
        printf "log_directory = '/data/logs/postgres'\n"
        printf "logging_collector = on\n"
        printf "log_filename = 'postgresql-%%Y-%%m-%%d.log'\n"
        printf "log_rotation_age = 1d\n"
        printf "log_rotation_size = 10MB\n"
    } >> "$PGDATA/postgresql.conf"

    # Start temp server to create database
    s6-setuidgid postgres pg_ctl -D "$PGDATA" -o "-c listen_addresses='localhost'" -w start -l /data/logs/postgres/init.log

    # Create the odoo database if it doesn't exist
    DB_NAME=$(jq -r '.odoo_db_name // "odoo"' "$OPTIONS_SOURCE")
    DB_EXISTS=$(s6-setuidgid postgres psql -U "$POSTGRES_USER" -tAc "SELECT 1 FROM pg_database WHERE datname='$DB_NAME';" 2>/dev/null || true)
    if [ "$DB_EXISTS" != "1" ]; then
        s6-setuidgid postgres createdb -U "$POSTGRES_USER" -O "$POSTGRES_USER" "$DB_NAME"
        echo "Database '$DB_NAME' created."
    fi

    # Stop temp server
    s6-setuidgid postgres pg_ctl -D "$PGDATA" -m fast -w stop

    echo "PostgreSQL initialization complete."
fi
